# deploy/systemd/units/justnews-boot-smoke.timer
[Unit]
Description=Run JustNews Boot Smoke Test after boot

[Timer]
OnBootSec=2min
Unit=justnews-boot-smoke.service

[Install]
WantedBy=timers.target
