# deploy/systemd/units/justnews-crawlers.target
[Unit]
Description=JustNews Crawlers Target
Wants=network-online.target postgresql.service mcp_bus.service
After=network-online.target postgresql.service mcp_bus.service

[Install]
WantedBy=multi-user.target
