# deploy/systemd/units/justnews-boot-smoke.service
[Unit]
Description=JustNews Boot Smoke Test
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
ExecStart=/home/adra/justnewsagent/JustNewsAgent/deploy/systemd/scripts/justnews-boot-smoke.sh
Nice=10
IOSchedulingClass=best-effort

[Install]
WantedBy=multi-user.target
